# lpie.service
[Unit]
Description=lpie reboot-rw service

[Service]
Type=oneshot
ExecStart=/opt/lpie/reboot-rw-svc

[Install]
WantedBy=multi-user.target
